# etc/systemd/system/prosody.service.d/override.conf
# Warte bis Netzwerk und Postgresql da sind
[Unit]
After=network-online.target postgresql.service
Requires=network-online.target

